# systemd/ai-assistant.service
[Unit]
Description=Ultahost AI Assistant Backend
After=network.target

[Service]
ExecStart=/usr/local/bin/ai-assistant
EnvironmentFile=/etc/ai-assistant.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
